# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

[program:vsftpd]
command=vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vsftpd.log
stderr_logfile=/var/log/supervisor/vsftpd.err
priority=1

[program:flask]
command=python -m flask run --host=0.0.0.0 --port=5000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/flask.log
stderr_logfile=/var/log/supervisor/flask.err
environment=
    FLASK_APP="app.py",
    FLASK_ENV="production",
    PYTHONUNBUFFERED="1",
    FTP_UPLOAD_DIR="/ftp-data/uploads",
    APP_PASSWORD="%(ENV_APP_PASSWORD)s",
    SECRET_KEY="%(ENV_SECRET_KEY)s",
    SMTP_HOST="%(ENV_SMTP_HOST)s",
    SMTP_PORT="%(ENV_SMTP_PORT)s",
    SMTP_USER="%(ENV_SMTP_USER)s",
    SMTP_PASS="%(ENV_SMTP_PASS)s",
    REPORT_EMAIL="%(ENV_REPORT_EMAIL)s"
priority=2
